# experiments/gmail-inbox/com.projectbrain.gmail-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.projectbrain.gmail-bridge</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>mkdir -p "$HOME/.project-brain/logs" &amp;&amp; exec "$HOME/.project-brain/app/venv/bin/python" "$HOME/.project-brain/app/experiments/gmail-inbox/bridge_v2.py" --apply &gt;&gt;"$HOME/.project-brain/logs/gmail-bridge.out.log" 2&gt;&gt;"$HOME/.project-brain/logs/gmail-bridge.err.log"</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PB_ALLOWED_SENDER</key>
    <string>trusted-sender@example.com</string>
  </dict>
  <key>StartInterval</key>
  <integer>60</integer>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
